# deploy/ansible/templates/x-postgres-backup.service.j2
[Unit]
Description=x-postgres-backup - PostgreSQL Backup Manager
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_dir }}
EnvironmentFile={{ app_dir }}/.env
ExecStart={{ app_dir }}/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port {{ app_port }}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
